# Scripts/crontab
MAILTO=""
00 03 * * 6 ${HOME}/bin/test_crash.sh
